# weather-station-web.service
[Unit]
Description=Weather Station Web Dashboard
After=network.target weather-station-api.service

[Service]
Type=simple
User=ethan
WorkingDirectory=/home/ethan/weather-station/apps/web
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/home/ethan/weather-station/.env

[Install]
WantedBy=multi-user.target

